# deployment/supervisor/microblog.conf
[program:t3]
command=/home/ubuntu/t3/venv/bin/gunicorn -b localhost:8000 -w 4 t3:app
directory=/home/ubuntu/t3
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
